# ubuntu/1c-gateway/systemd/1c-gateway.service
[Unit]
# Whitelist-прокси: единственная дверь мозга к 1С. Пропускает только read-методы,
# execute_code режет до 1С. Токен тулкита хранится тут (в EnvironmentFile), мозг его не знает.
Description=1C read-only gateway (execute_query allowlist)
After=network-online.target
Wants=network-online.target

[Service]
# EnvironmentFile: GW_TOOLKIT_TOKEN=<Bearer тулкита>, GW_UPSTREAM=http://192.168.56.1:6003, опц. GW_GATEWAY_TOKEN
EnvironmentFile=/etc/1c-gateway.env
ExecStart=/usr/bin/python3 /opt/1c-gateway/gateway.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
