# ubuntu/1c-etl/systemd/1c-etl.timer
[Unit]
Description=Ночная выгрузка 1С -> KB (03:00)

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
